# crontab
0 0 * * * mv /home/nate/*.temp /home/nate/.archive